# files/etc/uci-defaults/01_d
#!/bin/sh
uci set network.lan.ipaddr=192.168.2.1
uci commit network
uci set wireless.@wifi-device[0].disabled=0
uci set wireless.@wifi-iface[0].ssid=LEDE
uci commit wireless
uci commit
exit 0
